# scripts/crontab
@daily /home/chris/vocabulink/scripts/backup-database-to-s3.sh
