# etc/sv/postgresql/run
#!/bin/sh
exec chpst -upostgres:postgres /usr/bin/postmaster -D /var/lib/postgres/data 2>&1
