# etc/sv/vocabulink.cgi/run
#!/bin/sh
exec chpst -uhttp:http /home/jekor/vocabulink/vocabulink.cgi
