# etc/sv/vocabulink.cgi/run
#!/bin/sh
exec chpst -uwww-data:www-data /home/jekor/vocabulink/vocabulink.cgi 2>&1
